# systemd/pi-logger-gatherer.service
[Unit]
Description=pi-logger data gatherer

[Service]
Type=simple
ExecStart=/home/masneyb/data/pi-logger/bin/gather-data.sh
