# systemd/pi-logger-gatherer.timer
[Unit]
Description=pi-logger data gatherer timer

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=multi-user.target
