# systemd/pi-logger-create-graphs.timer
[Unit]
Description=pi-logger graph creator timer

[Timer]
OnUnitActiveSec=6h

[Install]
WantedBy=multi-user.target
