# systemd/pi-logger-gatherer.timer
[Unit]
Description=pi-logger data gatherer timer

[Timer]
OnCalendar=300s

[Install]
WantedBy=multi-user.target
